# monoverlay/factoryreset/data/cron/crontabs/root
*/3 * * * * /scripts/sensors_request 
0 1 * * * /scripts/shell_update_0to3600_delay >> /tmp/updatelog.tmp
* * * * * /scripts/update_cron_wlan0 
