# monoverlay/datafactoryreset/cron/crontabs/root
*/3 * * * * /scripts/sensors_request 
0 1 * * * /scripts/shell_update >> /tmp/updatelog.tmp
* * * * * /scripts/update_cron_wlan0 
